# deploy/systemd/robot-console-ai.service
[Unit]
Description=Matamoe Robot Console AI (Flask)
After=network-online.target
Wants=network-online.target

[Service]
User=matamoe
WorkingDirectory=/opt/robot/robot-console-ai
EnvironmentFile=-/opt/robot/robot-console-ai/.env
ExecStart=/opt/robot/robot-console-ai/.venv/bin/gunicorn -w 1 --threads 4 -b 0.0.0.0:8080 --timeout 180 --access-logfile - "app:APP"
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
