# deploy/systemd/telegram-robot-bot.service
[Unit]
Description=Telegram Robot Bot
After=network-online.target robot-console-ai.service
Wants=network-online.target

[Service]
Type=simple
User=matamoe
WorkingDirectory=/opt/robot/robot-console-ai
EnvironmentFile=-/opt/robot/robot-console-ai/.env
ExecStart=/opt/robot/robot-console-ai/.venv/bin/python /opt/robot/robot-console-ai/scripts/telegram_robot_bot.py
Restart=always
RestartSec=3
StandardOutput=append:/opt/robot/logs/telegram-robot-bot.log
StandardError=append:/opt/robot/logs/telegram-robot-bot.log

[Install]
WantedBy=multi-user.target
