# deploy/systemd/vlm-service.service
[Unit]
Description=Robot Console VLM Service
After=network-online.target
Wants=network-online.target

[Service]
User=matamoe
WorkingDirectory=/opt/robot/robot-console-ai
EnvironmentFile=-/opt/robot/robot-console-ai/.env
Environment=FLASK_APP=app_vlm.py
Environment=FLASK_RUN_PORT=8090
ExecStart=/opt/robot/robot-console-ai/.venv/bin/flask --app app_vlm.py run --host 0.0.0.0 --port 8090
Restart=always
RestartSec=2
StandardOutput=append:/opt/robot/logs/vlm-service.log
StandardError=append:/opt/robot/logs/vlm-service.log

[Install]
WantedBy=multi-user.target
